# robot/rovy.service
[Unit]
Description=ROVY Robot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rovy
WorkingDirectory=/home/rovy/rovy_client/robot
ExecStart=/bin/bash /home/rovy/rovy_client/robot/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=CAMERA_DEVICE=1
Environment=ROVER_SERIAL_DEVICE=/dev/ttyAMA0

[Install]
WantedBy=multi-user.target

